# root/etc/init.d/clawpanel
#!/bin/sh /etc/rc.common
# luci-app-clawpanel - init script for ClawPanel
#
# New directory scheme:
#   $CP_BASE_PATH/clawpanel/     ← 程序文件
#   $CP_BASE_PATH/configs/      ← 配置文件（clawpanel.json）
#   $CP_BASE_PATH/openclaw/     ← OpenClaw 数据
#   $CP_BASE_PATH/logs/         ← 日志文件

USE_PROCD=0
START=99
STOP=10

EXTRA_COMMANDS="setup status_service"
EXTRA_HELP="    setup            install ClawPanel\n    status_service   show service status"

_load_uci() {
    CP_BASE_PATH="$(uci -q get clawpanel.main.install_path || echo '')"
    if [ -z "$CP_BASE_PATH" ]; then
        echo "ClawPanel: install_path not set, run clawpanel-env setup first"
        return 1
    fi
    CP_INSTALL_PATH="${CP_BASE_PATH}/clawpanel"
    CP_CONFIG_PATH="${CP_BASE_PATH}/configs"
    CP_OPENCLAW_DIR="${CP_BASE_PATH}/openclaw"
    CP_LOGS_PATH="${CP_BASE_PATH}/logs"
    CP_BIN="${CP_INSTALL_PATH}/clawpanel"
    return 0
}

start_service() {
    _load_uci || return 1

    config_load clawpanel
    local enabled port
    config_get enabled main enabled "0"
    config_get port main port "19527"

    [ "$enabled" = "1" ] || {
        echo "ClawPanel: disabled in UCI (enabled=0)"
        return 0
    }

    [ -x "$CP_BIN" ] || {
        echo "ClawPanel: binary not found at $CP_BIN"
        return 1
    }

    port=$(echo "$port" | tr -cd '0-9')
    [ -z "$port" ] && port=19527

    # 确保目录存在
    mkdir -p "$CP_CONFIG_PATH" "$CP_OPENCLAW_DIR" "$CP_LOGS_PATH" 2>/dev/null

    # 首次安装时写入配置文件
    if [ ! -f "$CP_CONFIG_PATH/clawpanel.json" ]; then
        printf '{\n' > "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "port": %s,\n' "$port" >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "dataDir": "%s",\n' "$CP_OPENCLAW_DIR" >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "openClawDir": "%s",\n' "$CP_OPENCLAW_DIR" >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "openClawApp": "%s",\n' "/usr/local/lib/node_modules/openclaw" >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "openClawWork": "%s",\n' "$CP_OPENCLAW_DIR/work" >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "edition": "pro",\n' >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "logDir": "%s",\n' "$CP_LOGS_PATH" >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "jwtSecret": "clawpanel-secret-change-me",\n' >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "adminToken": "clawpanel",\n' >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '  "debug": false\n' >> "$CP_CONFIG_PATH/clawpanel.json"
        printf '}\n' >> "$CP_CONFIG_PATH/clawpanel.json"
    fi

    stop_service

    echo "Starting ClawPanel: $CP_BIN"
    echo "  install_path: $CP_INSTALL_PATH"
    echo "  config:       $CP_CONFIG_PATH/clawpanel.json"
    echo "  openclaw_dir: $CP_OPENCLAW_DIR"
    echo "  logs:         $CP_LOGS_PATH/"

    (
        export HOME="/root"
        export PATH="/usr/local/bin:$PATH"
        export NODE_ICU_DATA="/usr/local/share/icu"
        export LD_LIBRARY_PATH="/usr/local/lib:$LD_LIBRARY_PATH"
        setsid /bin/bash -c "$CP_BIN >> ${CP_LOGS_PATH}/clawpanel.log 2>&1 &"
    )

    local i=0
    while [ $i -lt 20 ]; do
        if netstat -tulnp 2>/dev/null | grep -q ":${port} "; then
            echo "ClawPanel started: port=$port"
            return 0
        fi
        sleep 1
        i=$((i + 1))
    done

    echo "ClawPanel: may have started, port $port not yet visible"
    return 0
}

stop_service() {
    _load_uci 2>/dev/null || true

    local port
    port=$(uci -q get clawpanel.main.port || echo "19527")
    port=$(echo "$port" | tr -cd '0-9')
    [ -z "$port" ] && port=19527

    for pid in $(pgrep -f "clawpanel" 2>/dev/null); do
        kill "$pid" 2>/dev/null
    done

    local i=0
    while [ $i -lt 10 ]; do
        if ! netstat -tulnp 2>/dev/null | grep -q ":${port} " 2>/dev/null; then
            break
        fi
        sleep 1
        i=$((i + 1))
    done
    return 0
}

reload_service() {
    stop_service
    sleep 1
    start_service
}

status_service() {
    _load_uci 2>/dev/null || true

    local port
    port=$(uci -q get clawpanel.main.port || echo "19527")
    local enabled
    enabled=$(uci -q get clawpanel.main.enabled || echo "0")

    echo "=== ClawPanel Status ==="
    echo "enabled:      $enabled"
    echo "install_path: $CP_BASE_PATH"
    echo "program:      $CP_INSTALL_PATH"
    echo "configs:      $CP_CONFIG_PATH"
    echo "openclaw_dir: $CP_OPENCLAW_DIR"
    echo "logs:         $CP_LOGS_PATH"
    echo "binary:       $CP_BIN"
    echo "port:         $port"
    if [ -x "$CP_BIN" ]; then
        echo "binary:       FOUND"
    else
        echo "binary:       NOT FOUND"
    fi
    if netstat -tulnp 2>/dev/null | grep -q ":${port} " 2>/dev/null; then
        echo "status:       RUNNING ✅"
    else
        echo "status:       NOT running"
    fi
}
