# root/etc/init.d/clawpanel
#!/bin/sh /etc/rc.common
#===========================================================
#  ClawPanel init.d — v1.3.0
#  注意: start/stop/status 已弃用，统一由 controller 管理
#  此脚本仅负责开机自启注册（enable/disable）
#===========================================================
USE_PROCD=0
START=99
STOP=10

start_service() {
    # 由 controller 的 action_service_ctl 统一管理
    # init.d 只做开机自启，不直接启动
    return 0
}

stop_service() {
    # 由 controller 统一杀进程
    return 0
}

reload_service() {
    return 0
}

status_service() {
    # 由 controller API 提供状态
    return 1
}
