# root/etc/init.d/clawpanel
#!/bin/sh /etc/rc.common
# luci-app-clawpanel - procd init script for ClawPanel

USE_PROCD=1
START=99
STOP=10

EXTRA_COMMANDS="setup status_service"
EXTRA_HELP="    setup            install ClawPanel\n    status_service   show service status"

# 閳光偓閳光偓 read UCI config (called each time start_service runs) 閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓閳光偓
_load_uci() {
    CP_BASE_PATH="$(uci -q get clawpanel.main.install_path || echo '')"
    if [ -z "$CP_BASE_PATH" ] || [ "${CP_BASE_PATH#/}" = "$CP_BASE_PATH" ]; then
        echo "ClawPanel: install_path not set, run clawpanel-env setup first"
        return 1
    fi
    CP_INSTALL_PATH="${CP_BASE_PATH}/clawpanel"
    CP_DATA="${CP_INSTALL_PATH}/data"
    CP_BIN="${CP_INSTALL_PATH}/clawpanel"
    CP_OPENCLAW_DIR="$(uci -q get clawpanel.main.openclaw_dir || echo '')"
    if [ -z "$CP_OPENCLAW_DIR" ] || [ "${CP_OPENCLAW_DIR#/}" = "$CP_OPENCLAW_DIR" ]; then
        CP_OPENCLAW_DIR="${CP_BASE_PATH}/.openclaw"
    fi
    return 0
}

start_service() {
    _load_uci || return 1

    local enabled port edition
    config_load clawpanel
    config_get enabled main enabled "0"
    config_get port main port "19527"
    config_get edition main edition "pro"

    [ "$enabled" = "1" ] || {
        echo "ClawPanel: disabled in UCI (enabled=0)"
        return 0
    }

    [ -x "$CP_BIN" ] || {
        echo "ClawPanel: binary not found at $CP_BIN, run clawpanel-env setup"
        return 1
    }

    mkdir -p "$CP_DATA" "$CP_OPENCLAW_DIR" 2>/dev/null

    # Ensure clawpanel.json has correct paths BEFORE ClawPanel starts.
    # ClawPanel may regenerate this file on first boot, so write it now.
    cat > "$CP_DATA/clawpanel.json" <<'JEOF'
{
  "port": ${port},
  "dataDir": "${CP_DATA}",
  "openClawDir": "${CP_OPENCLAW_DIR}",
  "openClawApp": "${CP_INSTALL_PATH}/app",
  "openClawWork": "${CP_INSTALL_PATH}/work",
  "edition": "pro",
  "jwtSecret": "clawpanel-secret-change-me",
  "adminToken": "clawpanel",
  "debug": false
}
JEOF

    port=$(echo "$port" | tr -cd '0-9')
    [ -z "$port" ] && port=19527

    # Kill old instances
    for pid in $(pgrep -f "clawpanel" 2>/dev/null); do
        kill "$pid" 2>/dev/null
    done

    # Wait for port to free
    local i=0
    while [ $i -lt 15 ]; do
        if ! netstat -tulnp 2>/dev/null | grep -q ":${port} "; then
            break
        fi
        usleep 200000 2>/dev/null || sleep 0.2
        i=$((i + 1))
    done

    # Read version from .version file
    local cp_ver
    if [ -f "${CP_INSTALL_PATH}/.version" ]; then
        cp_ver=$(cat "${CP_INSTALL_PATH}/.version" 2>/dev/null | tr -d '[:space:]')
    fi
    [ -z "$cp_ver" ] && cp_ver="5.3.3"
    # Strip edition prefix: pro-v5.3.3 閳?5.3.3
    cp_ver=$(echo "$cp_ver" | sed 's/^pro-v//;s/^lite-v//')

    echo "Starting ClawPanel: $CP_BIN --updater-standalone $cp_ver $CP_DATA $port $CP_OPENCLAW_DIR"

    procd_open_instance
    procd_set_param command "$CP_BIN" --updater-standalone "$cp_ver" "$CP_DATA" "$port" "$CP_OPENCLAW_DIR"
    procd_set_param respawn 3600 10 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param pidfile /var/run/clawpanel.pid
    procd_close_instance

    echo "ClawPanel started: port=$port openclaw_dir=$CP_OPENCLAW_DIR"
}

stop_service() {
    _load_uci 2>/dev/null || true

    local port
    port=$(uci -q get clawpanel.main.port || echo "19527")
    port=$(echo "$port" | tr -cd '0-9')
    [ -z "$port" ] && port=19527

    for pid in $(pgrep -f "clawpanel" 2>/dev/null); do
        kill "$pid" 2>/dev/null
    done

    local i=0
    while [ $i -lt 15 ]; do
        if ! netstat -tulnp 2>/dev/null | grep -q ":${port} "; then
            break
        fi
        usleep 200000 2>/dev/null || sleep 0.2
        i=$((i + 1))
    done
}

service_triggers() {
    procd_add_reload_trigger "clawpanel"
}

reload_service() {
    stop_service
    sleep 1
    start_service
}

status_service() {
    _load_uci 2>/dev/null || true

    local port
    port=$(uci -q get clawpanel.main.port || echo "19527")
    local enabled
    enabled=$(uci -q get clawpanel.main.enabled || echo "0")

    echo "=== ClawPanel Status ==="
    echo "enabled: $enabled"
    echo "install_path: $CP_BASE_PATH"
    echo "data_dir: $CP_DATA"
    echo "openclaw_dir: $CP_OPENCLAW_DIR"
    echo "binary: $CP_BIN"
    echo "port: $port"
    if [ -x "$CP_BIN" ]; then
        echo "binary_version: $("$CP_BIN" --version 2>/dev/null | head -1)"
    else
        echo "binary: not found"
    fi
    if netstat -tulnp 2>/dev/null | grep -q ":${port} "; then
        echo "status: RUNNING"
    else
        echo "status: NOT running"
    fi
}